# templates/grafana-backup.service.erb
[Unit]
Description=Run backup for grafana postgres

[Service]
Type=simple
ExecStart=/usr/local/bin/grafana-backup.sh <%= @datadir + '/backup' %>
<% if @postgres_watchdog -%>
ExecStartPost=/usr/bin/curl <%= @postgres_watchdog %>
<% end -%>
